# setup/supervisor/conf.d/celeryd.conf
[program:celery]
command=/home/ubuntu/Env/farmster-venv/bin/celery -E -A farmster_server worker -l info

directory=/home/ubuntu/farmster/farmster-server
numprocs=1
stdout_logfile=/home/ubuntu/farmster/farmster-server/logs/celery_worker.log
stderr_logfile=/home/ubuntu/farmster/farmster-server/logs/celery_worker.log
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true

; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=1000